# supervisord.conf
[supervisord]
nodaemon=true

[program:run.sh]
command=/bin/bash run.sh
directory=/
autorestart=false

[program:nginx]
command=nginx -g "daemon off";
autorestart=true
